# systemd/clockradio-mpd-autoconfig.service
[Unit]
Description=Auto-configure MPD audio_output for the current DAC
# Run before MPD so the daemon never sees stale hw:CARD= references
# pointing at a card the user has unplugged. ALSA enumerates USB
# audio devices during early boot via systemd-modules-load, so by
# the time we land here /proc/asound/cards is already populated.
After=systemd-modules-load.service
Before=mpd.service
ConditionPathExists=/etc/mpd.conf

[Service]
Type=oneshot
# Exit code 10 is the legacy "config rewritten" signal from the helper.
# The helper now handles the MPD restart itself when it detects MPD is
# already active (i.e. udev-triggered hot-plug path) so the systemd
# unit doesn't need an ExecStartPost — that path was brittle because
# $EXIT_STATUS isn't reliably exported to ExecStartPost across distros.
SuccessExitStatus=0 10
ExecStart=/usr/local/sbin/clockradio-mpd-autoconfig
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
